# root/etc/init.d/executable_swayosd-libinput-backend
#!/sbin/openrc-run

name=$RC_SVCNAME
description="SwayOSD libinput backend service"
command="/usr/bin/swayosd-libinput-backend"
command_user="root"
pidfile="/run/$RC_SVCNAME/$RC_SVCNAME.pid"
command_background="yes"

depend() {
	use logger
}

start_pre() {
        checkpath --directory --owner $command_user:$command_user --mode 0775 \
                /run/$RC_SVCNAME /var/log/$RC_SVCNAME
}
